# docker/crontab/cron.d/revengine.cron
# m h  dom mon dow   command
6 1 * * * /usr/src/app/packages/imports/bin/sailthru-s3-sync && /usr/src/app/packages/imports/bin/sailthru-import-yesterday
30 5 * * * cd /usr/src/app/ && /usr/local/bin/tsx /usr/src/app/packages/rfv/interactions.ts 
30 6 * * * cd /usr/src/app/ && /usr/local/bin/tsx /usr/src/app/packages/rfv/app_users.js